# root/etc/s6-overlay/s6-rc.d/init-pzserver-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

printf "Fixing permissons...\n"
    lsiown -R abc:abc /app
printf "Permissions fixed\n"